# inout-pten.timer
[Unit]
Description=Run PTEN Script Every 5 Minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
